# sqd-network-monitor.service
[Unit]
Description=SQD Node Monitoring Agent
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/sqd-network-monitor
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Make sure the config directory exists
ExecStartPre=/bin/mkdir -p /etc/sqd-network-monitor

[Install]
WantedBy=multi-user.target
